# rootfs/overlay/etc/systemd/system/rpmsg_json.service.d/override.conf
[Unit]
Wants=systemd-modules-load.service
After=systemd-modules-load.service

[Service]
ExecStartPre=/bin/sh -c 'for i in $(seq 1 30); do [ -e /sys/bus/platform/devices/78000000.r5f/remoteproc/remoteproc1/state ] && [ "$(cat /sys/bus/platform/devices/78000000.r5f/remoteproc/remoteproc1/state 2>/dev/null)" = "running" ] && [ -e /dev/rpmsg_ctrl1 ] && exit 0; sleep 1; done; exit 1'
Restart=on-failure
RestartSec=5
