# rootfs/overlay/etc/systemd/system/rpmsg_json.service.d/override.conf
[Unit]
Wants=systemd-modules-load.service
After=systemd-modules-load.service

[Service]
ExecStartPre=/bin/sh -c 'for i in $(seq 1 30); do for state in /sys/bus/platform/devices/78000000.r5f/remoteproc/remoteproc*/state; do [ -e "$state" ] || continue; [ "$(cat "$state" 2>/dev/null)" = "running" ] || continue; ls /dev/rpmsg_ctrl* >/dev/null 2>&1 && exit 0; done; sleep 1; done; exit 1'
Restart=on-failure
RestartSec=5
