# rootfs/overlays/sk-am64b-lab-r5f/etc/systemd/system/rpmsg_json.service.d/lab-mode.conf
[Unit]
ConditionPathExists=!/etc/ti-bringup/lab-r5f.mode
