# rootfs/overlays/sk-am64b-tsn-dscp-pcp/etc/systemd/system/ti-tsn-dscp-pcp-sk.service
[Unit]
Description=Re-apply SK TSN DSCP PCP network profile after boot
After=systemd-networkd.service multi-user.target
Wants=systemd-networkd.service

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 3
ExecStart=/bin/sh /usr/local/sbin/ti-tsn-dscp-pcp-sk-apply.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
